# v2/profile.d/aliases.sh
#!/bin/bash

alias ll='ls -alh'
alias c='clear'
alias v='vim'
